# launchd/com.michael.codex-as-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.michael.codex-as-api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/michael/.local/share/fnm/node-versions/v24.11.0/installation/bin/node</string>
        <string>/Users/michael/.dotfiles/bin/codex-as-api-serve.mjs</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/michael/Library/Logs/com.michael.codex-as-api.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/michael/Library/Logs/com.michael.codex-as-api.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/michael</string>
        <key>CODEX_HOME</key>
        <string>/Users/michael/.codex</string>
        <key>CODEX_AS_API_HOST</key>
        <string>127.0.0.1</string>
        <key>CODEX_AS_API_PORT</key>
        <string>18080</string>
        <key>NODE_PATH</key>
        <string>/Users/michael/.local/share/fnm/node-versions/v24.11.0/installation/lib/node_modules</string>
    </dict>
</dict>
</plist>
