# launchd/com.mshuffett.wacli-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mshuffett.wacli-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/wacli</string>
        <string>sync</string>
        <string>--follow</string>
    </array>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive / restart if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/michael/.wacli/sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/michael/.wacli/sync.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/michael</string>
    </dict>
</dict>
</plist>
