# launchd/com.michael.claude-evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.michael.claude-evening</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/michael/.dotfiles/bin/claude-evening-trigger</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/claude-evening.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-evening.log</string>
</dict>
</plist>
